# deploy/monitoring/aegis-external-probes.service
[Unit]
Description=AEGIS external federation probes

[Service]
Type=oneshot
ExecStart=/usr/local/bin/aegis-external-probes.sh
