# deploy/monitoring/aegis-prometheus-token.service
[Unit]
Description=Refresh Prometheus JWT bearer token for AEGIS /metrics

[Service]
Type=oneshot
ExecStart=/opt/aegis/monitoring/refresh-prometheus-token.sh
